# deploy/sim-server.service
# systemd unit for the native sim server, as an alternative to the container (docker-compose.yml).
# One process = one perpetual match (the Simulation class is per-match; run multiple units/ports
# or compose replicas for multiple concurrent matches — see docs/DEPLOY.md).
#
# Install:
#   dotnet publish server/SimServer.csproj -c Release -o /opt/wivuullegiance/sim
#   install -Dm640 deploy/sim-server.env.example /etc/wivuullegiance/sim.env   # then edit + set SIM_SECRET
#   install -Dm644 deploy/sim-server.service /etc/systemd/system/sim-server.service
#   systemctl daemon-reload && systemctl enable --now sim-server
[Unit]
Description=wivuullegiance native sim server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/wivuullegiance/sim
# Holds SIM_SECRET (required), STDB_HTTP / STDB_DB / STDB_TOKEN. Mode 0640, not world-readable.
EnvironmentFile=/etc/wivuullegiance/sim.env
ExecStart=/usr/bin/dotnet /opt/wivuullegiance/sim/SimServer.dll --port 8090
Restart=on-failure
RestartSec=2

# Hardening: run as a transient unprivileged user with a minimal view of the host.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
